# .sysconfig/runit/mikolaj-runsvdir/run
#!/bin/sh

exec 2>&1
exec chpst -P -u mikolaj runsvdir /home/mikolaj/services
