# .sysconfig/runit/runsvdir-mikolaj/run
#!/bin/sh

exec 2>&1
exec su - mikolaj -c 'runsvdir -P /home/mikolaj/services'
