# services/pulseaudio/finish
#!/bin/sh

../common.sh
exec ${PULSE} --kill
